# share/systemd/autograder%i.timer
[Unit]
Description=Daily run of the Autograder

[Timer]
OnCalendar=*-*-* 17:00:00
AccuracySec=30min

[Install]
WantedBy=multi-user.target

