# share/systemd/autograder%i.service
[Unit]
Description=Autograder Task


[Service]
ExecStart=/usr/bin/autograder %i
